# initscript
#! /bin/sh

### BEGIN INIT INFO
# Provides:          gstreamer-phone
# Should-Start:      
# Required-Start:    networking
# Required-Stop:     
# Default-Start:     2
# Default-Stop:      
# Short-Description: gstreamer-phone
# Description:       gstreamer-phone
### END INIT INFO

case "${1}" in
	start)
		/usr/src/gstreamer-phone/cronscript.py
		/usr/src/gstreamer-phone/gst-phone
		;;
esac

exit 0
